# deploy/arxiv-daily-ingest.timer
[Unit]
Description=Daily arXiv Paper Ingest Timer

[Timer]
OnCalendar=Sun,Mon,Tue,Wed,Thu *-*-* 20,21,22,23:30:00
Timezone=America/New_York
Persistent=true
Unit=arxiv-daily-ingest.service

[Install]
WantedBy=timers.target
