# deploy/arxiv-cleanup-tasks.service
[Unit]
Description=arXiv Task Queue Cleanup
After=network.target

[Service]
Type=oneshot
User=<user>
WorkingDirectory=<project_dir>
ExecStart=<project_dir>/.venv/bin/python3 scripts/cleanup_tasks.py
